# main.cron
0 0 * * * /Users/kalashnikova/Documents/универ/ВШЭ/2курс/1сем/DE/project/main.py